# deploy/systemd/db-backup.timer
[Unit]
Description=AI-AT-ADVENT daily DB backup timer
Documentation=https://github.com/marcoaduartemendes-source/AI-AT-ADVENT

[Timer]
# Runs once daily at 04:00 UTC — quiet hour for both US (midnight ET)
# and EU (06:00 CET), well clear of the orchestrator's 5-min ticks.
OnCalendar=*-*-* 04:00:00
# If the box was off/asleep at 04:00, fire on next boot.
Persistent=true
Unit=db-backup.service

[Install]
WantedBy=timers.target
