# deploy/systemd/scouts.service
[Unit]
Description=AI-AT-ADVENT scout sweep
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=aaa
Group=aaa
WorkingDirectory=/opt/ai-at-advent
# Source env via bash so multi-line values (KALSHI_PRIVATE_KEY) work.
# See orchestrator.service for the full reasoning.
TimeoutStartSec=240
ExecStart=/usr/bin/env bash -c 'set -a; . /etc/aaa.env; set +a; exec /opt/ai-at-advent/.venv/bin/python src/run_scouts.py --vacuum'

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/opt/ai-at-advent/data

Restart=on-failure
RestartSec=60
StartLimitBurst=3
StartLimitIntervalSec=1800
